# microdep-dev/etc/microdep/perfsonar-microdep-trace-ana.service
[Unit]
Description=perfSONAR Microdep traceroute analysis
After=network.target rabbitmq-server.service postgresql.service
Requires=rabbitmq-server.service postgresql.service

[Service]
Type=simple
ExecStart=/usr/lib/perfsonar/bin/microdep_commands/trace_event_reader.py -O /var/lib/logstash/microdep/trace-ana.json --pssrc amqp://localhost/?queue=traceroute --dbtype postgresql --geodb /etc/perfsonar/microdep/GeoLite2/GeoLite2-ASN.mmdb --namemap /etc/perfsonar/psconfig/pscheduler.d/toolkit-webui.json --topoevents

[Install]
WantedBy=multi-user.target
