# etc/microdep/perfsonar-microdep-restart.service
[Unit]
Description=Restart perfSONAR Microdep services

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart perfsonar-microdep-gap-ana.service

[Install]
WantedBy=multi-user.target