# microdep/perfsonar-microdep/scripts/perfsonar-microdep-watchconfig.path
[Unit]
Description=Watch for updates in toolkit-webui.json

[Path]
PathChanged=/etc/perfsonar/psconfig/pscheduler.d/toolkit-webui.json

